# _linux/databases.target
[Unit]
Description=A group of databases
Requires=dev_postgres.service
Before=dev_postgres.service

[Install]
RequiredBy=basic.target
