# _linux/queue.target
[Unit]
Description=A group of queues.
Requires=nats.service
Before=nats.service

[Install]
RequiredBy=basic.target
